# resources/gdal/supervisor/conf.d/cron.conf
[program:cron]
command = /usr/sbin/cron -f -L 15
user = root
autostart=true
autorestart=true