# resources/gdal/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log
autorestart=true