# resources/gdal/supervisor/conf.d/php5-fpm.conf
[program:php5-fpm]
command=/usr/sbin/php5-fpm -c /etc/php5/fpm
stdout_logfile=/var/log/supervisor/php5-fpm.log
stderr_logfile=/var/log/supervisor/php5-fpm-err.log
autorestart=true